# packages/occi/DEBIAN/postinst
#!/bin/sh
update-rc.d -f occi defaults
occi
